# incubator/condor-manager/images/condor-manager/container-files/etc/supervisord.d/htcondor.conf
[program:htcondor]
command=/usr/sbin/condor_master -f -t
