# incubator/condor-manager/images/condor-manager/container-files/etc/supervisord.d/tok-gen.conf
[program:generate_tokens]
command=/usr/local/bin/generate_tokens
startsecs=0
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
