# incubator/condor-login/images/condor-login/container-files/etc/supervisor.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -f /etc/ssh/sshd_config
