# incubator/gridftp/images/gridftp/supervisor/gridftp.conf
[program:gridftp]
command=/usr/local/sbin/gridftp-startup.sh
autorestart=true
startsecs=20
