# incubator/jupyter-notebook/images/jupyter-condor/container-files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
childlogdir = /var/log/supervisor

[program:sshd]
command=/usr/sbin/sshd -D

[program:jupyter]
command=start-notebook.sh --NotebookApp.password=%(ENV_JupyterPassword)s --NotebookApp.token=%(ENV_JupyterToken)s

[program:htcondor]
command=/usr/sbin/condor_master -f -t
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
