# stable/gridftp/images/gridftp/supervisor/gridftp.conf
[program:gridftp]
command=/usr/libexec/globus-gridftp-server-start
autorestart=true
startsecs=20
